# debian/green/green.service
[Service]
Type=oneshot
ExecStart=/usr/local/bin/green.sh
