# green.timer
[Unit]
Description=Ensure system is only up if being used

[Timer]
OnBootSec=1h
OnCalendar=minutely

[Install]
WantedBy=timers.target
