# green.service
[Service]
Type=oneshot
ExecStart=/home/nhooyr/src/nhooyr/dotfiles/linuxBin/green
