# xayah/green.service
[Service]
Type=oneshot
ExecStart=/home/nhooyr/src/nhooyr/dotfiles/xayah/green.js
